# deployment/templates/systemd/openclaw-task-manager-db.service
[Unit]
Description=OpenClaw Autonomous Task Manager (DB)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
KillMode=none
WorkingDirectory=%h/.openclaw/workspace
EnvironmentFile=-%h/.env
Environment=OPENCLAW_HOME=%h/.openclaw
ExecStartPre=/bin/bash %h/.openclaw/workspace/init-autonomous-tasks-db.sh
ExecStartPre=/bin/bash %h/.openclaw/workspace/migrate-task-statuses.sh
ExecStart=/bin/bash %h/.openclaw/workspace/autonomous-task-manager-db.sh
StandardOutput=append:%h/.openclaw/workspace/autonomous-task-manager-db.log
StandardError=append:%h/.openclaw/workspace/autonomous-task-manager-db.log

[Install]
WantedBy=default.target
