# deployment/templates/systemd/openclaw-telegram-commands.timer
[Unit]
Description=OpenClaw Telegram Task Commands Timer
Requires=openclaw-telegram-commands.service

[Timer]
OnBootSec=2min
OnUnitActiveSec=15s
Persistent=true
AccuracySec=5s

[Install]
WantedBy=timers.target
