# deployment/templates/systemd/twilio-call-monitor.timer
[Unit]
Description=Poll Twilio for new call recordings every 5 minutes

[Timer]
OnBootSec=3min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
